# kafka/supervisor/kafka-create-topic.conf
[program:create-kafka-topic]
command=/usr/bin/create-kafka-topic.sh
autostart=true
autorestart=false