# kafka/supervisor/kafka-create-topic.conf
[program:create-iot-topic]
command=/usr/bin/create-iot-topic.sh
autostart=true
autorestart=false